# files/aws-sync.timer
[Unit]
Description=Run AWS Reverse Proxy Config Daily Sync

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target